# poller/crontab
 * * * * * dotnet /app/poller.dll > /dev/stdout
