# poller/crontab
* * * * * echo "tick"
* * * * * dotnet /app/poller.dll
* * * * * ls /app
